# systemd/rtsp-archive.socket
[Socket]
ListenStream={{socket}}
SocketUser=http
SocketGroup=http
MaxConnections=4
FileDescriptorName=http.listen.socket

[Install]
RequiredBy={{name}}.service
