# systemd/rtsp-archive.service
[Unit]
Description={{description}}

[Service]
#Type=notify
ExecStart={{installdir}}/bin/{{name}}
TimeoutStartSec=20
#Restart=always
StandardOutput=syslog
User=nobody
Group=nobody

[Install]
WantedBy=multi-user.target
