# systemd/rtsp-archive.socket
[Socket]
ListenStream={{socket}}
MaxConnections=4
FileDescriptorName=http.listen.socket

[Install]
RequiredBy=rtsp-archive.service
