# systemd/rtsp-archive.service
[Unit]
Description={{description}}
After=network-online.target

[Service]
#Type=notify
ExecStart={{installdir}}/bin/{{name}} --config=/etc/{{name}}
TimeoutStartSec=20
#Restart=always
StandardOutput=syslog
User=nobody
Group=nobody
PrivateTmp=true
NoNewPrivileges=true
RuntimeDirectory={{name}}
StateDirectory={{name}}
ConfigurationDirectory={{name}}

[Install]
WantedBy=multi-user.target
